# systemd/media-stream-frontend.service
[Unit]
Description=Media Stream Frontend Server
After=network.target

[Service]
Type=simple
User=seithigal-youtube
WorkingDirectory=/home/seithigal-youtube/media_stream
Environment="PATH=/usr/bin:/usr/local/bin:/home/seithigal-youtube/.local/bin"
ExecStart=/usr/bin/npx serve -s dist -l 5173
Restart=always
RestartSec=10
StandardOutput=append:/home/seithigal-youtube/media_stream/logs/frontend.log
StandardError=append:/home/seithigal-youtube/media_stream/logs/frontend-error.log

[Install]
WantedBy=multi-user.target
